# setup/splash/bubuos-splash.service
# Replace YOUR_USER with your Linux username (e.g. pi)
[Unit]
Description=BubuOS Boot Splash
After=systemd-udev-settle.service
Before=bubuos.service

[Service]
Type=simple
User=root
Environment=SDL_VIDEODRIVER=x11
Environment=SDL_AUDIODRIVER=dummy
TTYPath=/dev/tty1
StandardInput=tty
ExecStart=/usr/bin/xinit /usr/bin/python3 /home/YOUR_USER/bubuos/setup/splash/bootsplash.py -- :9 vt1 -nocursor -nolisten tcp -keeptty
KillMode=control-group
TimeoutStopSec=3

[Install]
WantedBy=multi-user.target
