# setup/pocketos.service
[Unit]
Description=PocketOS Shell
After=multi-user.target NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/pocketos
Environment=SDL_VIDEODRIVER=kmsdrm
Environment=SDL_AUDIODRIVER=alsa
ExecStart=/usr/bin/python3 /home/pi/pocketos/main.py
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
